# deploy/telegram-crypto-market-crawler.service
[Unit]
Description=Telegram TRON Blacklist Crawler
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/telegram-crypto-market-bot
EnvironmentFile=/opt/telegram-crypto-market-bot/.env
ExecStart=/usr/bin/npm run crawl:blacklist
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
